# card/var/lib/dpkg/info/dnscrypt-proxy.list
/.
/etc
/etc/dnscrypt-proxy
/etc/dnscrypt-proxy/dnscrypt-proxy.toml
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/dnscrypt-proxy-resolvconf.service
/lib/systemd/system/dnscrypt-proxy.service
/lib/systemd/system/dnscrypt-proxy.socket
/usr
/usr/sbin
/usr/sbin/dnscrypt-proxy
/usr/share
/usr/share/doc
/usr/share/doc/dnscrypt-proxy
/usr/share/doc/dnscrypt-proxy/README.md
/usr/share/doc/dnscrypt-proxy/changelog.Debian.armhf.gz
/usr/share/doc/dnscrypt-proxy/changelog.Debian.gz
/usr/share/doc/dnscrypt-proxy/changelog.gz
/usr/share/doc/dnscrypt-proxy/copyright
/usr/share/doc/dnscrypt-proxy/examples
/usr/share/doc/dnscrypt-proxy/examples/example-blacklist.txt
/usr/share/doc/dnscrypt-proxy/examples/example-cloaking-rules.txt
/usr/share/doc/dnscrypt-proxy/examples/example-dnscrypt-proxy.toml
/usr/share/doc/dnscrypt-proxy/examples/example-forwarding-rules.txt
