# card/var/lib/dpkg/info/setserial.conffiles
/etc/init.d/etc-setserial
/etc/init.d/setserial
