# card/var/lib/dpkg/info/setserial.list
/.
/bin
/bin/setserial
/etc
/etc/init.d
/etc/init.d/etc-setserial
/etc/init.d/setserial
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/etc-setserial.service
/lib/systemd/system/setserial.service
/usr
/usr/share
/usr/share/doc
/usr/share/doc/setserial
/usr/share/doc/setserial/README
/usr/share/doc/setserial/README.Debian.gz
/usr/share/doc/setserial/changelog.Debian.gz
/usr/share/doc/setserial/copyright
/usr/share/doc/setserial/pcmcia.repair
/usr/share/doc/setserial/serial.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/setserial
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/setserial.8.gz
/var
/var/lib
/var/lib/setserial
